# launchd/ai.openclaw.gateway.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Comment</key>
	<string>OpenClaw Gateway (v2026.3.24)</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>ANTHROPIC_API_KEY</key>
		<string>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</string>
		<key>DISCORD_BOT_TOKEN</key>
		<string>MTQ4MzE5MzIzMjcwNjc2ODk5OA.GyOY9k.mHBnIPM5jA9-wB5zzxFd8m-Fkz2e5aXBdsVsXM</string>
		<key>GITHUB_TOKEN</key>
		<string>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</string>
		<key>HOME</key>
		<string>/Volumes/BotCentral/Users/milo</string>
		<key>NODE_EXTRA_CA_CERTS</key>
		<string>/etc/ssl/cert.pem</string>
		<key>NODE_USE_SYSTEM_CA</key>
		<string>1</string>
		<key>NVIDIA_API_KEY</key>
		<string>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</string>
		<key>NVIDIA_NIM_API_KEY</key>
		<string>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</string>
		<key>OPENCLAW_GATEWAY_PORT</key>
		<string>18789</string>
		<key>OPENCLAW_LAUNCHD_LABEL</key>
		<string>ai.openclaw.gateway</string>
		<key>OPENCLAW_SERVICE_KIND</key>
		<string>gateway</string>
		<key>OPENCLAW_SERVICE_MARKER</key>
		<string>openclaw</string>
		<key>OPENCLAW_SERVICE_VERSION</key>
		<string>2026.3.24</string>
		<key>OPENCLAW_SYSTEMD_UNIT</key>
		<string>openclaw-gateway.service</string>
		<key>OPENCLAW_WINDOWS_TASK_NAME</key>
		<string>OpenClaw Gateway</string>
		<key>OPENROUTER_API_KEY</key>
		<string>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</string>
		<key>PATH</key>
		<string>/opt/homebrew/opt/node/bin:/Volumes/BotCentral/Users/milo/.local/bin:/Volumes/BotCentral/Users/milo/.npm-global/bin:/Volumes/BotCentral/Users/milo/bin:/Volumes/BotCentral/Users/milo/.volta/bin:/Volumes/BotCentral/Users/milo/.asdf/shims:/Volumes/BotCentral/Users/milo/.bun/bin:/Volumes/BotCentral/Users/milo/Library/Application Support/fnm/aliases/default/bin:/Volumes/BotCentral/Users/milo/.fnm/aliases/default/bin:/Volumes/BotCentral/Users/milo/Library/pnpm:/Volumes/BotCentral/Users/milo/.local/share/pnpm:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
		<key>PERPLEXITY_API_KEY</key>
		<string>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</string>
		<key>TELEGRAM_BOT_TOKEN</key>
		<string>8774680410:AAFoqOTynJZUPvPmlBWDkSqFixsAzmfa60c</string>
		<key>TMPDIR</key>
		<string>/var/folders/4g/85lybl0s3kd9v5yt1khy9vmr0000gn/T/</string>
	</dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>ai.openclaw.gateway</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/opt/node/bin/node</string>
		<string>/opt/homebrew/lib/node_modules/openclaw/dist/index.js</string>
		<string>gateway</string>
		<string>--port</string>
		<string>18789</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Volumes/BotCentral/Users/milo/.openclaw/logs/gateway.err.log</string>
	<key>StandardOutPath</key>
	<string>/Volumes/BotCentral/Users/milo/.openclaw/logs/gateway.log</string>
	<key>ThrottleInterval</key>
	<integer>1</integer>
	<key>Umask</key>
	<integer>63</integer>
</dict>
</plist>
