# launchd/ai.openclaw.gateway.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Comment</key>
	<string>OpenClaw Gateway (v2026.3.24)</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>ANTHROPIC_API_KEY</key>
		<string>${ANTHROPIC_API_KEY}</string>
		<key>DISCORD_BOT_TOKEN</key>
		<string>${DISCORD_BOT_TOKEN}</string>
		<key>GITHUB_TOKEN</key>
		<string>${GITHUB_TOKEN}</string>
		<key>HOME</key>
		<string>/Volumes/BotCentral/Users/milo</string>
		<key>NODE_EXTRA_CA_CERTS</key>
		<string>/etc/ssl/cert.pem</string>
		<key>NODE_USE_SYSTEM_CA</key>
		<string>1</string>
		<key>NVIDIA_API_KEY</key>
		<string>${NVIDIA_API_KEY}</string>
		<key>NVIDIA_NIM_API_KEY</key>
		<string>${NVIDIA_NIM_API_KEY}</string>
		<key>OPENCLAW_GATEWAY_PORT</key>
		<string>18789</string>
		<key>OPENCLAW_LAUNCHD_LABEL</key>
		<string>ai.openclaw.gateway</string>
		<key>OPENCLAW_SERVICE_KIND</key>
		<string>gateway</string>
		<key>OPENCLAW_SERVICE_MARKER</key>
		<string>openclaw</string>
		<key>OPENCLAW_SERVICE_VERSION</key>
		<string>2026.3.24</string>
		<key>OPENCLAW_SYSTEMD_UNIT</key>
		<string>openclaw-gateway.service</string>
		<key>OPENCLAW_WINDOWS_TASK_NAME</key>
		<string>OpenClaw Gateway</string>
		<key>OPENROUTER_API_KEY</key>
		<string>${OPENROUTER_API_KEY}</string>
		<key>PATH</key>
		<string>/opt/homebrew/opt/node/bin:/Volumes/BotCentral/Users/milo/.local/bin:/Volumes/BotCentral/Users/milo/.npm-global/bin:/Volumes/BotCentral/Users/milo/bin:/Volumes/BotCentral/Users/milo/.volta/bin:/Volumes/BotCentral/Users/milo/.asdf/shims:/Volumes/BotCentral/Users/milo/.bun/bin:/Volumes/BotCentral/Users/milo/Library/Application Support/fnm/aliases/default/bin:/Volumes/BotCentral/Users/milo/.fnm/aliases/default/bin:/Volumes/BotCentral/Users/milo/Library/pnpm:/Volumes/BotCentral/Users/milo/.local/share/pnpm:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
		<key>PERPLEXITY_API_KEY</key>
		<string>${PERPLEXITY_API_KEY}</string>
		<key>TELEGRAM_BOT_TOKEN</key>
		<string>${TELEGRAM_BOT_TOKEN}</string>
		<key>TMPDIR</key>
		<string>/var/folders/4g/85lybl0s3kd9v5yt1khy9vmr0000gn/T/</string>
	</dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>ai.openclaw.gateway</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/opt/node/bin/node</string>
		<string>/opt/homebrew/lib/node_modules/openclaw/dist/index.js</string>
		<string>gateway</string>
		<string>--port</string>
		<string>18789</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Volumes/BotCentral/Users/milo/.openclaw/logs/gateway.err.log</string>
	<key>StandardOutPath</key>
	<string>/Volumes/BotCentral/Users/milo/.openclaw/logs/gateway.log</string>
	<key>ThrottleInterval</key>
	<integer>1</integer>
	<key>Umask</key>
	<integer>63</integer>
</dict>
</plist>
